# crontab
* * * * * root php /var/www/studip/cli/cronjob-worker.php >> /dev/null 2>&1
